# container-data/etc/services.d/cron/run
#!/usr/bin/execlineb -P
with-contenv
s6-envuidgid nginx
foreground {
    cd /app
    artisan schedule:run
}
sleep 60
