# container-data/etc/services.d/nginx/run
#!/usr/bin/execlineb -P
/usr/sbin/nginx -c /etc/nginx/nginx.conf -g "daemon off;"
