# container-data/etc/services.d/php7-fpm/run
#!/usr/bin/execlineb -P
with-contenv
/usr/sbin/php-fpm7 --nodaemonize --fpm-config /etc/php7/php-fpm.conf
